# conf/crontab
MAILTO=tauberer@govtrack.us
14 3 * * * { cd /home/dotcloud/current; ../env/bin/python update_meeting_times.py; }
