# conf/crontab
MAILTO=tauberer+anc@govtrack.us
14 3 * * * { cd /home/ancbrigade/web; python conf/cronjob.py; ./fcgi; }
