# conf/crontab
MAILTO=tauberer+anc@govtrack.us
14 3 * * * { cd /home/ancbrigade/web; LANG=en_US.UTF-8 python conf/cronjob.py; USER=ancbrigade ./fcgi > /dev/null; }
